# tests/crontab
# Clears the swarm every night at midnight
0 0 * * * sudo k-swarm-prune
